# init.d/bitlbee
#!/sbin/runscript

description="Brings instant messaging (XMPP, MSN, Yahoo!, AIM, ICQ, Twitter) to IRC"
pidfile=/run/bitlbee/bitlbee.pid
command=bitlbee
command_args="$BITLBEE_OPTS"

depend() {
   need logger net
}
