# init.d/login
#!/sbin/runscript

description="Remove /var/run/nologin created by systemd tmpfiles.d"

depend() {
	after tmpfilesd.boot
}

start() {
	if [[ -f /var/run/nologin ]]; then
		rm -f /var/run/nologin
		return $?
	else
		return 0
	fi
}
